# images/consul/consul_servers.service
[Unit]
Description=consul_servers

[Service]
ExecStart=/bin/sh -c "/usr/bin/consul_servers.sh > /etc/consul-servers"
